# docker/crontab
# WB-Parser cron schedule: twice daily at 06:00 and 18:00
0 6,18 * * * www-data cd /var/www/html && php artisan wb:fetch:all >> /var/log/cron.log 2>&1
# Empty line required at end of crontab
