# docker/crontab
0 6,18 * * * cd /var/www/html && php artisan wb:fetch:all >> /var/log/cron.log 2>&1
